# GeecsCAGateway/deploy/geecs-ca-gateway.service
# GEECS → EPICS Channel Access gateway — systemd unit for the lab deployment
# on abmx-Super-Server (192.168.6.14), alongside the GEECS MySQL DB and Tiled.
#
# Install (from the repo checkout on the box):
#   sudo cp deploy/geecs-ca-gateway.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now geecs-ca-gateway
#
# Upgrade / DB-resync (a restart IS both — config rebuilds live from the DB):
#   git pull && poetry install && sudo systemctl restart geecs-ca-gateway
#
# Adjust User/paths for a different box; see deploy/README.md and DEPLOYMENT.md.

[Unit]
Description=GEECS to EPICS Channel Access gateway (Undulator)
# The gateway reads the GEECS DB at startup. If MySQL runs on this box under a
# different unit name (mysqld/mariadb), fix the name; if the DB is remote,
# drop it — Restart=on-failure retries until the DB is reachable either way.
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=abmx
# Poetry resolves its cache/venvs under $HOME; systemd does not always set it.
Environment=HOME=/home/abmx
WorkingDirectory=/home/abmx/GEECS-Plugins/GeecsCAGateway
# Serve CA only on the lab interface (write-safety residual — no CA-level auth).
Environment=EPICS_CAS_INTF_ADDR_LIST=192.168.6.14
Environment=EPICS_CAS_BEACON_ADDR_LIST=192.168.6.255
ExecStart=/home/abmx/.local/bin/poetry run python -m geecs_ca_gateway --experiment Undulator
Restart=on-failure
RestartSec=10
# Remote restart: a CA put to <Experiment>:CAGateway:RESTART makes the gateway
# exit with code 86, which systemd treats as "relaunch" (clean, not a failure).
# This is also how the served set resyncs after a GEECS DB edit — no SSH needed.
SuccessExitStatus=86
RestartForceExitStatus=86
SyslogIdentifier=geecs-ca-gateway

[Install]
WantedBy=multi-user.target
